# systemd/nebiros-server.service
[Unit]
Description = Nebiros Web Server

[Service]
ExecStart=/opt/app/Nebiros/web/NebirosServer

[Install]
WantedBy=multi.user.target
